# systemd/syncthing@root.service.d/override.conf
[Service]
InaccessibleDirectories=/etc /home
CapabilityBoundingSet=~CAP_SYS_ADMIN
